# supervisor.newyorkfieldguide.conf
[program:newyorkfieldguide]
command=gunicorn_django -c gunicorn.py
directory=/srv/newyorkfieldguide
user=sclay
autostart=true
autorestart=true
#redirect_stderr=True
priority=991
stopsignal=HUP
